# config/systemd/user/skwd-wall-color.service
[Unit]
Description=Generate waybar colors from skwd-wall

[Service]
Type=oneshot
ExecStart=%h/.dot/config/waybar/scripts/generate-colors.sh
